# .config/systemd/user/lf-handler.service
[Unit]
Description=LF D-Bus handler
After=default.target

[Service]
ExecStart=%h/.config/lf/utils/file-handler

[Install]
WantedBy=default.target
